# spider.1688/crontab/crawl_all.sh
#!/usr/bin/env bash
cd /home/apps/jd_spider
sh crontab/crawl_a.sh &
sh crontab/crawl_b.sh &
sh crontab/crawl_c.sh &
sh crontab/crawl_d.sh &
sh crontab/crawl_e.sh &
sh crontab/crawl_f.sh &
sh crontab/crawl_g.sh &
sh crontab/crawl_h.sh &
sh crontab/crawl_i.sh &
sh crontab/crawl_j.sh &
sh crontab/crawl_k.sh &
sh crontab/crawl_l.sh &
sh crontab/crawl_m.sh &
sh crontab/crawl_n.sh &
sh crontab/crawl_o.sh &



