# spider.bjsat/crontab/check_proxy_to_bjsat.sh
#!/usr/bin/env bash
cd /home/apps/bjsat.gov.cn.spider
python proxy/check_proxy_to_bjsat.py &

