# spider.qichacha_app/crontab/crawl.sh
#!/usr/bin/env bash
scrapy crawl qichacha_app_spider