# spider.baic/crontab/check_proxy_to_baic.sh
#!/usr/bin/env bash
cd /home/apps/spider.baic
python proxy/check_proxy_to_baic.py &

